# rpi-wifi-led.service
[Unit]
Description=Raspberry Pi WiFi LED control server
After=network.target

[Service]
# Adjust User and paths to match where you cloned the repo.
User=pi
WorkingDirectory=/home/pi/rpi-wifi-led
# Uses the system python3 with the apt-installed flask/gpiozero/lgpio packages.
ExecStart=/usr/bin/python3 /home/pi/rpi-wifi-led/app.py
Restart=on-failure
RestartSec=3
Environment=LED_PINS=18,23
Environment=PORT=5000

[Install]
WantedBy=multi-user.target
